# config/cloudflared.service
[Unit]
Description=Cloudflare Tunnel
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel run
EnvironmentFile=/etc/cloudflared/env
Restart=always
RestartSec=5
User=cloudflared
Group=cloudflared

[Install]
WantedBy=multi-user.target
