# config/litellm.service
[Unit]
Description=LiteLLM Proxy
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
ExecStart=/usr/local/bin/litellm --config /etc/litellm/config.yaml --host 127.0.0.1 --port 4000
EnvironmentFile=/etc/litellm/env
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
User=litellm
Group=litellm
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ReadWritePaths=/var/lib/litellm /home/litellm/.cache
MemoryMax=1536M
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
